# update-motd.d/00-header
#!/bin/bash

. /etc/os-release
. /etc/orangepi-release

printf "\nWelcome to %s (%s %s %s)\n\n" "$PRETTY_NAME $KERNELID" "$(uname -o)" "$(uname -r)" "$(uname -m)"
